# scripts/com.sovereignsignal.portforward.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.sovereignsignal.portforward</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/kubectl</string>
        <string>-n</string>
        <string>sovereign-signal</string>
        <string>port-forward</string>
        <string>service/sovereign-signal</string>
        <string>8080:80</string>
    </array>
    <key>KeepAlive</key>
    <true/>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/andrewboss/Library/Logs/sovereign-signal-portforward.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/andrewboss/Library/Logs/sovereign-signal-portforward-error.log</string>
</dict>
</plist>
